# root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

mkdir -p /var/log/nginx/
exec php-fpm;
